# modules/test/host/test-done/mkosi.extra/usr/local/lib/systemd/system/test-done.service
[Unit]
Description=Barrier that triggers test services
DefaultDependencies=no

Wants=snphost-config-commit.service
After=snphost-config-commit.service

[Service]
Type=oneshot
ExecStart=/usr/bin/true
RemainAfterExit=yes

[Install]
WantedBy=test.target
