# modules/build/host/targets/mkosi.extra/usr/local/lib/systemd/system/test.target
# Target reached when sev-verify tests are complete.
# Job services should use: WantedBy=test.target

[Unit]
Description=test stage

After=system.target
Requires=system.target

Wants=test-done.service
After=test-done.service
